# upstart/MaliitServer-launcher.conf.in
# Copyright (c) 2017 LG Electronics, Inc.
#
# Starts MaliitServer when
#   1) Certain conditions are met during the boot-up. (static mode only)
#   2) The MaliitServer was stopped abnormally. This actually respawns it in error case.
#   3) Someone emits ime-activate signal. (effective in dynamic mode)
#   4) Requested explicitly, i.e. start MaliitServer-launcher

start on (minimal-boot-done and settingsservice-ready and lsm-ready) or (ime-activate) or (stopped MaliitServer)

script
    PMLOG_CONTEXT="IME"
    PMLOG_KV="file=\"MaliitServer-launcher.conf\""
    PmLogCtl def $PMLOG_CONTEXT > /dev/null 2>&1

    RESPAWN=false

    if [ "$UPSTART_EVENTS" == "stopped" -a "$JOB" == "MaliitServer" ]; then
        if [ "$RESULT" == "ok" ]; then
            PmLogCtl logkv $PMLOG_CONTEXT warning upstart $PMLOG_KV "MaliitServer terminated normally, stop respawning"
            exit 0
        elif [ "$RESULT" == "failed" ]; then
            PmLogCtl logkv $PMLOG_CONTEXT warning upstart $PMLOG_KV "MaliitServer terminated abnormally, respawn in a few seconds (EXIT_STATUS=$EXIT_STATUS, EXIT_SIGNAL=$EXIT_SIGNAL)"
            sleep 5
            RESPAWN=true
        fi

        # Check wayland connection
        if [ -z "$WAYLAND_DISPLAY" ]; then
            WAYLAND_DISPLAY="wayland-0"
        fi
        if ( ! netstat -lx | grep -q "$XDG_RUNTIME_DIR/$WAYLAND_DISPLAY" ); then
            # Compositor seems down, don't rush.
            PmLogCtl logkv $PMLOG_CONTEXT warning upstart $PMLOG_KV "Wait for the wayland connection ready..."
            sleep 30
        fi
    fi

    STATIC=$(luna-send -n 1 -a com.webos.service.ime luna://com.webos.service.config/getConfigs '{"configNames":["com.webos.service.ime.static"]}' | grep -oE "\"com.webos.service.ime.static\"[ \t]*:[ \t]*(true|false)" | grep -oE "true|false" || echo "true")

    if [ "$STATIC" == true ]; then
        PmLogCtl logkv $PMLOG_CONTEXT info upstart $PMLOG_KV "Starting up MaliitServer (static mode)"
        start MaliitServer STATIC=$STATIC RESPAWN=$RESPAWN
    else
        if [ "$UPSTART_EVENTS" == "ime-activate" -o "$RESPAWN" == "true" ]; then
            PmLogCtl logkv $PMLOG_CONTEXT info upstart $PMLOG_KV "Starting up MaliitServer (dynamic mode, UPSTART_EVENTS=$UPSTART_EVENTS, RESPAWN=$RESPAWN)"
            start MaliitServer STATIC=$STATIC RESPAWN=$RESPAWN
        else
            PmLogCtl logkv $PMLOG_CONTEXT warning upstart $PMLOG_KV "No start MaliitServer (dynamic mode, UPSTART_EVENTS=$UPSTART_EVENTS)"
        fi
    fi
end script
